# supervisord.conf
[supervisord]
nodaemon=true

[program:1]
command=/usr/bin/v2ray run -config ./config.json > /dev/null 2>&1
autorestart=true
autostart=true

